# plugins/general-plugins/skills/loop_macos-workspace/iteration-1/eval-0-daily-shell/without_skill/outputs/output.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
    "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>

    <!-- Unique reverse-DNS label — launchd uses this as the job identity -->
    <key>Label</key>
    <string>com.user.claude.sync-anthropic-skills</string>

    <!-- Full path to the script; ~ is NOT expanded by launchd, use absolute path -->
    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/jhonsmith/.claude/my-marketplace/bin/sync-anthropic-skills.sh</string>
    </array>

    <!-- Daily trigger: fire at 09:00 every day -->
    <key>StartCalendarInterval</key>
    <dict>
        <key>Hour</key>
        <integer>9</integer>
        <key>Minute</key>
        <integer>0</integer>
    </dict>

    <!-- Redirect stdout so the run is auditable -->
    <key>StandardOutPath</key>
    <string>/Users/jhonsmith/.claude/my-marketplace/logs/sync-anthropic-skills.log</string>

    <!-- Redirect stderr separately so errors are easy to grep -->
    <key>StandardErrorPath</key>
    <string>/Users/jhonsmith/.claude/my-marketplace/logs/sync-anthropic-skills.error.log</string>

    <!-- Run even if the previous interval was missed (e.g. laptop was asleep) -->
    <key>RunAtLoad</key>
    <false/>

</dict>
</plist>
