# plugins/general-plugins/skills/loop_macos-workspace/iteration-2/eval-1-30min-downloads/with_skill/outputs/output.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.loop-macos.check-if-there-are-any-new</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>find /Users/jhonsmith/Downloads -maxdepth 1 -mmin -30 -not -name '.*' -newer /Users/jhonsmith/Downloads | while read -r f; do echo "$(date '+%Y-%m-%d %H:%M:%S') NEW: $f"; done</string>
    </array>

    <key>StartInterval</key>
    <integer>1800</integer>

    <key>StandardOutPath</key>
    <string>/Users/jhonsmith/Library/Logs/loop-macos/com.user.loop-macos.check-if-there-are-any-new.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jhonsmith/Library/Logs/loop-macos/com.user.loop-macos.check-if-there-are-any-new.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
