# plugins/general-plugins/skills/loop_macos-workspace/iteration-1/eval-2-2h-python/with_skill/outputs/output.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.user.loop-macos.python-scripts-cleanup-tmp-py</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>-lc</string>
        <string>python ~/scripts/cleanup_tmp.py</string>
    </array>

    <key>StartInterval</key>
    <integer>7200</integer>

    <key>StandardOutPath</key>
    <string>/Users/jhonsmith/Library/Logs/loop-macos/com.user.loop-macos.python-scripts-cleanup-tmp-py.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/jhonsmith/Library/Logs/loop-macos/com.user.loop-macos.python-scripts-cleanup-tmp-py.log</string>

    <key>RunAtLoad</key>
    <false/>
</dict>
</plist>
